# app/DEBIAN/postinst
#!/usr/bin/env bash

# GUI only on desktop.
if type gnome-shell &> /dev/null; then
    # Package dbus-x11 is needed for GUI root scripts like kz-install.
    apt-get install --yes dbus-x11 zenity
fi

# Snap install mainly for Debian.
if ! snap list core &> /dev/null; then
    systemctl enable --now snapd.socket
    systemctl restart snapd.seeded.service
    snap install core
fi

true
