# DEBIAN/postinst
#!/bin/bash

# Voor Debain.
snap install core
systemctl enable --now snapd.socket

true
